# aptura-os/packages/aptura-desktop/etc/skel/.config/autostart/aptura-welcome.desktop
[Desktop Entry]
Type=Application
Name=Aptura Welcome
Comment=Show first-run status and setup shortcuts
Exec=aptura-welcome
Icon=aptura-welcome
OnlyShowIn=Aptura;
Terminal=false
X-GNOME-Autostart-enabled=true
