# aptura-os/config/live-build/config/includes.chroot/etc/skel/.config/autostart/aptura-welcome.desktop
[Desktop Entry]
Type=Application
Name=Aptura Welcome
Comment=Start the Aptura Flow welcome experience
Exec=aptura-flow --welcome
OnlyShowIn=GNOME;
X-GNOME-Autostart-enabled=true
