# aptura-os/packages/aptura-branding/etc/profile.d/aptura-branding.sh
export APTURA_BRAND_NAME="Aptura OS"
export APTURA_BRAND_EDITION="Classic XFCE"
export APTURA_BRAND_ACCENT="#000080"
