# aptura-os/packages/aptura-branding/etc/profile.d/aptura-branding.sh
export APTURA_BRAND_NAME="Aptura OS"
export APTURA_BRAND_EDITION="Shell"
export APTURA_BRAND_ACCENT="#00d9c0"
